# config/x86_64-192.168.0.14/hello-openwrt/config/rc.local
# 读写权限
mount -o remount rw /
# 设置tmp缓存空间为128M
mount -o remount,size=128M tmpfs /tmp
# 挂载Swap缓存
swapon '/www/swapfile'
# 禁用ICMP
echo 1 > /proc/sys/net/ipv4/icmp_echo_ignore_all
echo 1 > /proc/sys/net/ipv6/icmp/echo_ignore_all
# 30分钟后，关闭SSH 和 HTTP 及tyyd
(sleep 1800 && /etc/init.d/dropbear stop && /etc/init.d/uhttpd stop && /etc/init.d/ttyd stop) &
exit 0
